# src/server-Lenovo-T15/systemd/wg-quick@.service.d/10-t15.conf
# /etc/systemd/system/wg-quick@.service.d/10-t15.conf — managed by t15
#
# Applies to all wg-quick@<name> instances (template drop-in). Binds the tunnel
# to t15.target so stopping the target stops the tunnel, and auto-restart on
# failure so a transient kernel-module hiccup recovers without the watchdog.
#
# The cycle prevention is handled on the dnsmasq side now (bind-dynamic +
# no Requires=wg-quick@wg0 in dnsmasq's drop-in), so we don't need to touch
# wg-quick's default After=/Wants= here.

[Unit]
PartOf=t15.target

[Service]
Restart=on-failure
RestartSec=15s
