# src/server-Lenovo-T15p/systemd/create-ap0.service
[Unit]
Description=Create ap0 virtual WiFi interface for IoT access point
Before=hostapd.service
After=sys-subsystem-net-devices-wlp0s20f3.device
Wants=sys-subsystem-net-devices-wlp0s20f3.device

[Service]
Type=oneshot
RemainAfterExit=yes
# Create virtual AP interface on the same physical radio as wlp0s20f3
ExecStart=/sbin/iw dev wlp0s20f3 interface add ap0 type __ap
# Assign static IP for the IoT LAN gateway
ExecStart=/sbin/ip addr add 192.168.10.1/24 dev ap0
ExecStart=/sbin/ip link set ap0 up
ExecStop=/sbin/iw dev ap0 del

[Install]
WantedBy=multi-user.target
