# src/server-Lenovo-T15p/systemd/t15p-watchdog.timer
[Unit]
Description=T15p VPN Server Watchdog Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
