# src/server-Lenovo-T15/systemd/dnsmasq.service.d/10-t15.conf
# /etc/systemd/system/dnsmasq.service.d/10-t15.conf — managed by t15
#
# dnsmasq picks up interfaces dynamically (bind-dynamic in iot.conf), so it
# can start before wg0 exists and bind to it when wg-quick@wg0 brings it up.
# We only require t15-lan (enp0s31f6 has the IP dnsmasq listens on for DHCP).
#
# Deliberately NOT Requires=/After= wg-quick@wg0 — that combination created a
# boot-time cycle via dnsmasq's default Before=nss-lookup.target:
#   dnsmasq → nss-lookup.target → wg-quick@wg0 → dnsmasq

[Unit]
Requires=t15-lan.service
After=t15-lan.service
PartOf=t15.target

[Service]
Restart=on-failure
RestartSec=10s
