# src/server-Lenovo-T15/systemd/nftables.service.d/10-t15.conf
# /etc/systemd/system/nftables.service.d/10-t15.conf — managed by t15
#
# Restart the firewall on failure. NOT PartOf=t15.target: the firewall must
# stay up even when the VPN is stopped, to keep unsolicited inbound dropped.

[Service]
Restart=on-failure
RestartSec=10s
